# crontab
# m h dom mon dow command
*/5 *   * * *   ~/PiTempSensor/bin/logtemp
