# crontab
# m h dom mon dow user  command
*/5 *   * * *   root    /home/pi/bin/write_stats_to_db.sh
